# airootfs/etc/xdg/autostart/equestria-os-tutorial.desktop
[Desktop Entry]
Type=Application
Name=Equestria OS Tour
Exec=equestria-os-tutorial
Icon=equestria-os-tutorial
X-KDE-Autostart-after=panel